# packaging/sensord_command.socket
[Unit]
Description=Sensor command socket

[Socket]
SocketUser=sensor
SocketGroup=input
ListenStream=/tmp/sensord_command_socket
SocketMode=0777
PassCredentials=yes
Accept=false
SmackLabelIPIn=*
SmackLabelIPOut=@
Service=sensord.service
